# web-inject/web-inject-server/_system/etc/systemd/system/juice-update.service
[Unit]
After=juice-shop.service
[Service]
Type=simple
ExecStart=/home/ubuntu/update.sh

[Install]
WantedBy=multi-user.target
Alias=juice-update

