# web-sde/web-sde-server/_system/etc/systemd/system/juice-hashcopy.service
[Service]
Type=simple
ExecStart=/sde_configuration.sh

[Install]
WantedBy=multi-user.target
Alias=juice-hashcopy

